# systemd/mirror-world-dailyschedulereminder.service
[Unit]
Description=Mirror World - DailyScheduleReminder (schedule reminders + optional ticket startup messages)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=rsadmin
WorkingDirectory=/home/rsadmin/bots/mirror-world
ExecStart=/bin/bash /home/rsadmin/bots/mirror-world/RSAdminBot/run_bot.sh dailyschedulereminder
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
